# upstart.conf
#!upstart

description "node-slides"
author      "wavded"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
env NODE_ENV=production

script
   exec /usr/local/bin/node /var/www/node-slides/app.js >> /var/log/node/node-slides.log 2>&1
end script
